# red/custom-daemon/go-api-daemon/go-api.service
[Unit]
Description=Go REST API
[Service]
User=root
WorkingDirectory=/root/tools/
ExecStart=/root/tools/log4shell-white-box/red/custom-daemon/go-api-daemon/run-go-api.sh

SuccessExitStatus=143
TimeoutStopSec=10
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
