# red/custom-daemon/java-http-daemon/java-http.service
[Unit]
Description=Java HTTP Service
[Service]
User=cookie
WorkingDirectory=/home/cookie/tools/
ExecStart=/home/cookie/tools/run-java-http.sh

SuccessExitStatus=143
TimeoutStopSec=10
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
